# deploy/systemd/ddag-logical-backup.timer
[Unit]
Description=Daily DDAG encrypted logical backup

[Timer]
OnCalendar=*-*-* 02:15:00
Persistent=true
RandomizedDelaySec=10m
Unit=ddag-logical-backup.service

[Install]
WantedBy=timers.target
