# deploy/systemd/ddag-backup-metrics.timer
[Unit]
Description=Refresh DDAG backup and PostgreSQL metrics every minute

[Timer]
OnBootSec=2m
OnUnitActiveSec=1m
Persistent=true
RandomizedDelaySec=10s
Unit=ddag-backup-metrics.service

[Install]
WantedBy=timers.target
